# xdg/sv/fluxbox/run
#!/bin/bash

exec /usr/bin/fluxbox -display ${DISPLAY} -log ${HOME}/.local/log/fluxbox.log
